# rootfs/rootfs/etc/profile.d/boot2docker.sh
#!/bin/sh

echo "boot2docker: $(cat /etc/version)"
echo "             $(cat /etc/boot2docker)"
