# rootfs/rootfs/etc/profile.d/boot2docker.sh
#!/bin/sh

echo "Boot2Docker version $(cat /etc/version), build $(cat /etc/boot2docker)"
docker -v # "Docker version 1.3.0-dev, build ba14ddf-dirty"
